# spotweb/rootfs/etc/periodic/hourly/spotweb
#!/usr/bin/with-contenv bashio
# ==============================================================================
# Runs the hourly task to update spots
# ==============================================================================

/usr/bin/php8 /app/retrieve.php
